# roles/ovp_ir/templates/etc/systemd/system/ovpn_peer_rotate.service
[Unit]
Description=The script to change the Peer IP address

[Service]
ExecStart=/bin/bash {{ ovpn_root }}/peer_rotate.sh

[Install]
WantedBy=default.target

