# roles/ovp_ir/templates/etc/systemd/system/ovpn_peer_rotate_stunnel.service
[Unit]
Description=The script to change the Peer IP address for stunnel

[Service]
ExecStart=/bin/bash {{ ovpn_root }}/peer_rotate_stunnel.sh

[Install]
WantedBy=default.target

