# roles/ovp_ir/templates/etc/systemd/system/stunnel-ovpn-server.service
[Unit]
Description=Stunnel Service for OpenVPN

[Service]
Type=forking
WorkingDirectory=/opt/ovpn/
ExecStart=/usr/bin/stunnel4 /etc/stunnel/stunnel-ovpn-server.conf
SuccessExitStatus=0
Restart=always

[Install]
WantedBy=multi-user.target
